# PKGBUILD/mpd_oled/mpd_oled.service
[Unit]
Description=MPD OLED Display

[Service]
User=root
EnvironmentFile=/etc/default/mpd_oled
ExecStart=/usr/bin/mpd_oled $OPTS
KillMode=control-group
KillSignal=SIGTERM
ExecStopPost=/bin/sh -c '/usr/bin/rm -f /tmp/cava*'
 
[Install]
WantedBy=multi-user.target
