# PKGBUILD/mpd_oled/mpd_oled.service
[Unit]
Description=MPD OLED Display

[Service]
Environment=HOME=/tmp
EnvironmentFile=/etc/default/mpd_oled
ExecStart=/usr/bin/mpd_oled $OPTS

[Install]
WantedBy=multi-user.target
